# deploy/tetramem-api.service
[Unit]
Description=TetraMem-XL v7.0 API Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/workspace/sunorm-space-memory
ExecStart=/usr/bin/python3 start_api_v2.py
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

Environment=TETRAMEM_STORAGE=/root/.openclaw/workspace/sunorm-space-memory/tetramem_data_v2
Environment=TETRAMEM_PORT=8000
Environment=TETRAMEM_HOST=127.0.0.1
Environment=TETRAMEM_UI_PASSWORD=Hive@2026
Environment=TETRAMEM_RESOLUTION=5
Environment=TETRAMEM_CHECKPOINT_INTERVAL=120

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tetramem-api

WatchdogSec=60

[Install]
WantedBy=multi-user.target
